# applications/luci-app-amule/root/etc/uci-defaults/40_luci-amule
#!/bin/sh

mkdir -p /etc/amule
touch /etc/amule/addresses.dat
touch /etc/amule/ipfilter_static.dat
touch /etc/amule/shareddir.dat

exit 0
