# applications/luci-app-rsyncd/root/etc/uci-defaults/luci-rsyncd
#!/bin/sh

chmod 600 /etc/rsyncd.secrets
chmod 600 /etc/rsyncd.password

exit 0
