# applications/luci-app-php7/root/etc/uci-defaults/40_luci-php7
#!/bin/sh

/etc/init.d/rpcd reload
exit 0
