# applications/luci-app-dockerman/root/etc/uci-defaults/40_luci-dockerman
#!/bin/sh

/etc/init.d/dockerd disable
/etc/init.d/dockerd stop
exit 0
