# applications/luci-app-webd/root/etc/uci-defaults/40_luci-webd
#!/bin/sh

mkdir -p /etc/webd/.Trash

exit 0
