# applications/luci-app-mosdns/root/etc/uci-defaults/40_luci-mosdns
#!/bin/sh

mkdir -p /etc/mosdns/rule
touch /etc/mosdns/rule/whitelist.txt
touch /etc/mosdns/rule/blocklist.txt

exit 0
