# applications/luci-app-statistics/root/etc/uci-defaults/40_luci-statistics
#!/bin/sh

# symlink for busybox httpd
[ -x /usr/sbin/httpd ] && [ ! -h /www/rrdimg ] && \
	ln -s /tmp/rrdimg /www/rrdimg

exit 0
