# applications/luci-app-haproxy-tcp/root/etc/uci-defaults/luci-haproxy
#!/bin/sh

/etc/init.d/haproxy disable
/etc/init.d/haproxy stop

exit 0
